# git/dotfiles/system/etc/rules.d/60-ssd-scheduler.rules
# Set deadline scheduler for SSDs
ACTION=="add|change", SUBSYSTEM=="block", KERNEL=="nvme*", RUN+="/bin/sh -c 'echo kyber > /sys/block/nvme0n1/queue/scheduler"
